# bybot-telegram.service
[Unit]
Description=BYBOT Telegram Bot (Integrated)
After=network.target

[Service]
Type=simple
User=mikevance
WorkingDirectory=/home/mikevance/bots/bybot
ExecStart=/home/mikevance/bots/venv/bin/python3.12 main.py
Restart=always

[Install]
WantedBy=multi-user.target 